# launchd/com.iia.us-rankings.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.iia.us-rankings</string>

    <!-- caffeinate -is: prevent idle + system sleep during execution -->
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/caffeinate</string>
        <string>-is</string>
        <string>/Users/edward.song/.local/bin/uv</string>
        <string>run</string>
        <string>/Users/edward.song/Desktop/Stock/scripts/fetch_rankings.py</string>
        <string>--market</string>
        <string>US</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/edward.song/Desktop/Stock</string>

    <!--
        兩個時間點觸發，腳本內部用 zoneinfo 判斷 NYSE 是否已收盤：
        - 04:30 台灣時間 = 夏令 16:30 ET（NYSE 已收盤） → 抓取成功
        - 06:00 台灣時間 = 冬令 17:00 ET（NYSE 已收盤）→ 抓取成功
        夏令時：04:30 抓取後 DB 有資料，06:00 被 skip-if-exists 略過
        冬令時：04:30 NYSE 未收盤被跳過，06:00 正常抓取

        Weekday 偏移：美股週一收盤後台灣已是週二凌晨
        → Weekday 2(週二) ~ 6(週六) 對應美股週一到週五
    -->
    <key>StartCalendarInterval</key>
    <array>
        <!-- 04:30 台灣時間，週二~週六（對應美股週一~週五） -->
        <dict><key>Hour</key><integer>4</integer><key>Minute</key><integer>30</integer><key>Weekday</key><integer>2</integer></dict>
        <dict><key>Hour</key><integer>4</integer><key>Minute</key><integer>30</integer><key>Weekday</key><integer>3</integer></dict>
        <dict><key>Hour</key><integer>4</integer><key>Minute</key><integer>30</integer><key>Weekday</key><integer>4</integer></dict>
        <dict><key>Hour</key><integer>4</integer><key>Minute</key><integer>30</integer><key>Weekday</key><integer>5</integer></dict>
        <dict><key>Hour</key><integer>4</integer><key>Minute</key><integer>30</integer><key>Weekday</key><integer>6</integer></dict>
        <!-- 06:00 台灣時間，週二~週六（冬令 fallback） -->
        <dict><key>Hour</key><integer>6</integer><key>Minute</key><integer>0</integer><key>Weekday</key><integer>2</integer></dict>
        <dict><key>Hour</key><integer>6</integer><key>Minute</key><integer>0</integer><key>Weekday</key><integer>3</integer></dict>
        <dict><key>Hour</key><integer>6</integer><key>Minute</key><integer>0</integer><key>Weekday</key><integer>4</integer></dict>
        <dict><key>Hour</key><integer>6</integer><key>Minute</key><integer>0</integer><key>Weekday</key><integer>5</integer></dict>
        <dict><key>Hour</key><integer>6</integer><key>Minute</key><integer>0</integer><key>Weekday</key><integer>6</integer></dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/edward.song/Desktop/Stock/logs/us_rankings.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/edward.song/Desktop/Stock/logs/us_rankings_err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/edward.song/.local/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/edward.song</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
