# launchd/com.iia.catchup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.iia.catchup</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/caffeinate</string>
        <string>-is</string>
        <string>/Users/edward.song/.local/bin/uv</string>
        <string>run</string>
        <string>/Users/edward.song/Desktop/Stock/scripts/catchup.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/edward.song/Desktop/Stock</string>

    <!-- Trigger 1: run immediately when agent is loaded (boot / launchctl load) -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Trigger 2: run every 30 minutes as a safety net -->
    <key>StartInterval</key>
    <integer>1800</integer>

    <!-- Trigger 3: run when network connects (resolv.conf changes on DHCP/VPN) -->
    <key>WatchPaths</key>
    <array>
        <string>/private/var/run/resolv.conf</string>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/edward.song/Desktop/Stock/logs/catchup.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/edward.song/Desktop/Stock/logs/catchup.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/edward.song/.local/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/edward.song</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>
</dict>
</plist>
